# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:application_server]
command=gunicorn --bind 0.0.0.0:8000 -w 1 inventory_management.wsgi:application
directory=/project/
user=PrioritySoft
autostart=true
autorestart=true
redirect_stderr=true
